# service/PIRSensorTester/run
#!/bin/sh
exec /home/pi/PIRSensorTester 2>&1
